# apps/disher-backend-3.0/ops/disher-backup.service
# disher Postgres off-box backup — oneshot unit driven by disher-backup.timer.
#
# Install (on the VPS, as root):
#   1. Edit the two paths below (WorkingDirectory + the script path) to match
#      where the repo is checked out, and set the env values.
#   2. cp ops/disher-backup.{service,timer} /etc/systemd/system/
#   3. systemctl daemon-reload && systemctl enable --now disher-backup.timer
#   4. Verify: systemctl start disher-backup.service ; journalctl -u disher-backup
#
# Runs the backup as root so it can read a root-owned rclone.conf. pg-backup.sh
# already reads POSTGRES_* from .env.production in WorkingDirectory; the offsite +
# heartbeat knobs are passed here.
[Unit]
Description=disher Postgres off-box backup (pg_dump -> integrity -> offsite -> heartbeat)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/disher/apps/disher-backend-3.0
Environment=BACKUP_DIR=/srv/disher/backups
# Offsite + heartbeat are OFF: both need real credentials we do not have yet.
# Leaving a placeholder value here is worse than leaving it unset — pg-backup.sh
# fails the whole run when `rclone copyto` fails, so a fake remote = zero backups.
# To enable offsite: install rclone, write a root-owned /root/.config/rclone/rclone.conf,
# then uncomment and point this at the real remote.
#Environment=RCLONE_REMOTE=b2:disher-backups
# To enable heartbeat: create a check on healthchecks.io / Better Stack and paste its
# ping URL here. Unset = no alert if this timer silently dies.
#Environment=HEALTHCHECK_URL=https://hc-ping.com/<uuid>
ExecStart=/usr/bin/env bash scripts/pg-backup.sh
Nice=10
IOSchedulingClass=idle
